# systemd/resol2mqtt.service
[Unit]
Description=resol2mqtt bridge daemon
Documentation=https://github.com/resol2mqtt/resol2mqtt
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/resol2mqtt
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
Alias=resol2mqtt.service
